# systemd/llmrouter.service
[Unit]
Description=LLMRouter local LiteLLM proxy
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/llmrouter
EnvironmentFile=%h/llmrouter/.env
ExecStart=%h/llmrouter/start-router.sh
ExecStop=/usr/bin/env LLMROUTER_SKIP_SYSTEMD_STOP=1 %h/llmrouter/stop-router.sh

[Install]
WantedBy=default.target
